# systemd/backup.timer
[Unit]
Description=Backup
Wants=network-online.target

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
