# systemd/backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=backup -c /etc/backup/backup.conf create

[Install]
WantedBy=timers.target
