# systemd/backup.service
[Unit]
Description=Backup
OnFailure=discord-notification@%n.service

[Service]
Type=oneshot
ExecStart=backup

[Install]
WantedBy=timers.target