# systemd/backup-notification@.service
[Unit]
Description=Send backup notifications

[Service]
Type=oneshot
ExecStart=notify-ntfy %I
