# systemd/discord-notification@.service
[Unit]
Description=Send discord notification regarding service %I

[Service]
Type=oneshot
ExecStart=notify-discord %I
